# .systemd/arca.service
[Unit]
Description=openhands backend for arca
After=atlas.service
Requires=atlas.service

[Service]
Type=simple
ExecStartPre=/bin/sleep 5
ExecStart=/app/.venv/bin/uvicorn openhands.server.listen:app --host 0.0.0.0 --port 3000
Restart=always
RestartSec=5
User=arca
Group=arca
WorkingDirectory=/app
Environment=/app/.venv/bin

[Install]
WantedBy=multi-user.target